# supervisord.conf
[supervisord]
; Roda em foreground (o container vive enquanto este processo viver)
nodaemon=true
; Em containers costuma-se usar root; os processos abaixo podem rodar como www-data
user=root
; Opcional: registra o estado do supervisor (útil para troubleshooting)
; logfile=/dev/stdout
; logfile_maxbytes=0

; ----------------------------- PHP-FPM ---------------------------------------
[program:php-fpm]
; -F => foreground (não daemoniza); supervisor controla o ciclo de vida
command=php-fpm -F
; Sobe junto com o container e reinicia em caso de crash
autostart=true
autorestart=true
; Sinal para parada elegante
stopsignal=QUIT
; Mantém hierarquia de processos para kill em cascata
stopasgroup=true
killasgroup=true
; Dá alguns segundos para subir antes de considerar "online"
startsecs=2
; Prioridade menor: sobe primeiro (quanto menor, mais cedo)
priority=5
; Logs para stdout/stderr do container (sem rotação local)
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Usuário efetivo do processo (combina com permissões do PHP)
user=www-data

; ------------------------------ NGINX ----------------------------------------
[program:nginx]
; Força a usar o nginx.conf fornecido; o arquivo já tem "daemon off;"
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
autostart=true
autorestart=true
stopsignal=QUIT
stopasgroup=true
killasgroup=true
startsecs=2
; Nginx depende do php-fpm estar de pé primeiro (pequena janela)
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
; Usuário do master process; workers herdam do nginx.conf (user www-data;)
user=root
